# launchd/com.nutridesi.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nutridesi.server</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/caffeinate</string>
    <string>-dims</string>
    <string>/usr/local/bin/node</string>
    <string>server.js</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/swapnil/Documents/Whatsapp calorie bot/Swapnils first App</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/swapnil/Library/Logs/nutridesi.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/swapnil/Library/Logs/nutridesi.log</string>
</dict>
</plist>
