# launchd/com.nutridesi.ngrok.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nutridesi.ngrok</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/ngrok</string>
    <string>http</string>
    <string>3000</string>
    <string>--log=stdout</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/swapnil/Library/Logs/nutridesi-ngrok.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/swapnil/Library/Logs/nutridesi-ngrok.log</string>
</dict>
</plist>
